# eeee.conf
[program:eeee]
command=python3 /home/curator2/eeee/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/eeee/test.err.log
stdout_logfile=/home/curator2/eeee/test.out.log